# supervisor.conf
[supervisord]
nodaemon=true
[program:hexo]
command= hexo generate --watch &
[program:nginx]
command= nginx -g "daemon off;"
