# supervisor.conf
[supervisord]
nodaemon=true
[program:hexo]
command= cd /blog && hexo server -i 0.0.0.0